# yazi-fm.service
[Unit]
Description=Yazi File Manager IPC Service
Documentation=https://github.com/user/yazi-frontend
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStartPre=/usr/bin/rm -rf %h/.cache/quickshell/qmlcache
ExecStart=/usr/bin/qs -c yazi-fm
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
